# tiktok-bot.service
[Unit]
Description=Telegram TikTok Downloader Bot
Documentation=https://github.com/ragilmalik/telegram-tiktok-downloader
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=ragilmalik
WorkingDirectory=/home/ragilmalik/telegram-tiktok-downloader
ExecStart=/usr/bin/node /home/ragilmalik/telegram-tiktok-downloader/bot.js

# CRITICAL: Load environment variables from .env file
EnvironmentFile=/home/ragilmalik/telegram-tiktok-downloader/.env

# Additional environment
Environment=NODE_ENV=production

# Restart policy
Restart=on-failure
RestartSec=10

# Logging - systemd journal (use journalctl to view logs)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tiktok-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
