# tiktok-bot.service
[Unit]
Description=Telegram TikTok Downloader Bot
Documentation=https://github.com/ragilmalik/telegram-tiktok-downloader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/telegram-tiktok-downloader
ExecStart=/opt/node22/bin/node /home/user/telegram-tiktok-downloader/bot.js

# CRITICAL: Load environment variables from .env file
EnvironmentFile=/home/user/telegram-tiktok-downloader/.env

# Additional environment
Environment=NODE_ENV=production

# Restart policy with maximum retry limit to prevent infinite restarts
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

# Logging - systemd journal (use journalctl to view logs)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tiktok-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
